# configs/lemonbar/config/lemonbar-scheduler.service
[Unit]
Description=Lemonbar scheduler service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/my-lemonbar-scheduler
