# configs/lemonbar/config/lemonbar-clock.timer
[Unit]
Description=Lemonbar clock update timer

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
AccuracySec=500ms

[Install]
WantedBy=timers.target
