# configs/lemonbar/config/lemonbar-clock.service
[Unit]
Description=My clock updater

[Service]
Type=oneshot
ExecStart=%h/.scripts/my-lemonbar-update
