# configs/lemonbar/config/lemonbar-scheduler.timer
[Unit]
Description=Lemonbar scheduler

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
AccuracySec=100ms

[Install]
WantedBy=timers.target
