# configs/lemonbar/config/lemonbar.service
[Unit]
Description=My lemonbar main runner

[Service]
ExecStart=%h/.local/bin/my-lemonbar-main

[Install]
WantedBy=default.target
