# configs/lemon-modules/config/lemon-modules-scheduler.timer
[Unit]
Description=Lemon modules scheduler timer.

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
