# configs/lemon-modules/config/lemon-modules.service
[Unit]
Description=Lemon modules main runner.

[Service]
ExecStart=%h/.local/bin/lemon-modules-server

[Install]
WantedBy=default.target
