# configs/lemon-modules/config/lemon-modules-scheduler.service
[Unit]
Description=Lemon modules scheduler service.

[Service]
Type=oneshot
ExecStart=%h/.local/bin/lemon-modules-scheduler
